# infra/production/systemd/balikha-backup.service
[Unit]
Description=Balikha database backup (pg_dump -> private R2)
After=network-online.target postgresql.service
Wants=network-online.target
OnFailure=balikha-job-failure-alert@%n.service

[Service]
Type=oneshot
EnvironmentFile=/etc/balikha/backup.env
ExecStart=/opt/balikha/current/infra/production/backup.sh
