# infra/production/systemd/balikha-job-failure-alert@.service
[Unit]
Description=Balikha scheduled-job failure alert for %I
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=balikha-app
Group=balikha-app
SupplementaryGroups=systemd-journal
WorkingDirectory=/opt/balikha/current
EnvironmentFile=/etc/balikha/production.env
Environment=NODE_ENV=production
Environment=HOME=/var/lib/balikha
ExecStart=/opt/balikha/current/node_modules/.bin/tsx db/scripts/notify-job-failure.ts %I
