# infra/production/systemd/balikha-job-failure-alert@.service
[Unit]
# %i (raw instance), NOT %I: systemd unescaping maps '-' -> '/', which corrupts
# hyphenated unit names (balikha-weekly-digest -> balikha/weekly/digest). The
# OnFailure=...@%n.service idiom makes the instance the failing unit's literal
# name (e.g. balikha-weekly-digest.service), which %i returns verbatim.
Description=Balikha scheduled-job failure alert for %i
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=balikha-app
Group=balikha-app
SupplementaryGroups=systemd-journal
WorkingDirectory=/opt/balikha/current
EnvironmentFile=/etc/balikha/production.env
Environment=NODE_ENV=production
Environment=HOME=/var/lib/balikha
ExecStart=/opt/balikha/current/node_modules/.bin/tsx db/scripts/notify-job-failure.ts %i
