# infra/production/systemd/balikha-bootstrap-admin.service
[Unit]
Description=Balikha static-admin bootstrap (one-shot; ADMIN_EMAIL/ADMIN_PASSWORD from production.env)
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=oneshot
User=balikha-app
Group=balikha-app
WorkingDirectory=/opt/balikha/current
EnvironmentFile=/etc/balikha/production.env
Environment=NODE_ENV=production
Environment=HOME=/var/lib/balikha
ExecStart=/opt/balikha/current/node_modules/.bin/tsx db/scripts/bootstrap-admin.ts
