# infra/production/systemd/balikha-orders-tick.timer
[Unit]
Description=Run the Balikha order sweep hourly

[Timer]
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=120

[Install]
WantedBy=timers.target
